# tag-osx/Library/LaunchAgents/com.boneskull.claude-watcher.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnvironmentVariables</key>
	<dict>
		<key>PATH</key>
		<string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
	</dict>
	<key>KeepAlive</key>
	<true/>
	<key>Label</key>
	<string>com.boneskull.claude-watcher</string>
	<key>ProgramArguments</key>
	<array>
		<string>/opt/homebrew/opt/node@24/bin/node</string>
		<string>/Users/boneskull/.claude/plugins/cache/boneskull-plugins/claude-watcher/0.1.2/dist/bin/claude-watcher.js</string>
		<string>daemon</string>
	</array>
	<key>RunAtLoad</key>
	<true/>
	<key>StandardErrorPath</key>
	<string>/Users/boneskull/.config/claude-watcher/logs/daemon.err</string>
	<key>StandardOutPath</key>
	<string>/Users/boneskull/.config/claude-watcher/logs/daemon.log</string>
</dict>
</plist>
